# supervisord.conf
[supervisord]
nodaemon=true

[program:emqttd]
command=/opt/emqttd/bin/emqttd start